# scripts/systemd/nxbtd.service
[Unit]
Description=PlayGate NXBT Bluetooth daemon
Documentation=https://github.com/playgate/playgate
# nuxbt must start only after bluetoothd is up and running.
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
# nxbtd.py path is set by install-host.sh; default: /opt/playgate/nxbt-daemon/nxbtd.py
ExecStart=/usr/bin/python3 /opt/playgate/nxbt-daemon/nxbtd.py
WorkingDirectory=/opt/playgate/nxbt-daemon

# nuxbt requires root (or CAP_NET_ADMIN) to talk to bluetoothd via DBus
# and to register the controller's HID profile.
User=root

# Automatic restart on unexpected exit
Restart=on-failure
RestartSec=5

# Create /run/nxbt at service start (owned by root, world-readable so
# the playgate user can connect to the UNIX socket).
RuntimeDirectory=nxbt
RuntimeDirectoryMode=0755

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nxbtd

[Install]
WantedBy=multi-user.target
